# deploy/whoscored/systemd/whoscored-daily-issuer.service
[Unit]
Description=Freeze and sign one bounded WhoScored rollout-wave approval
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
EnvironmentFile=/etc/data-platform/whoscored-daily-issuer.env
LoadCredential=approval-hmac:/etc/data-platform/credentials/whoscored-approval-hmac
LoadCredential=owner-hmac:/etc/data-platform/credentials/whoscored-owner-hmac
LoadCredential=issuance-ledger-hmac:/etc/data-platform/credentials/whoscored-issuance-ledger-hmac
ExecStart=/usr/local/libexec/whoscored-daily-issuer
TimeoutStartSec=30min
RuntimeDirectory=whoscored-daily-issuer
RuntimeDirectoryMode=0700
NoNewPrivileges=true
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE
AmbientCapabilities=
PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
ProtectSystem=strict
ReadWritePaths=/run/whoscored-daily-issuer
ProtectHome=read-only
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
